# systemd/osrt-pkglistgen@.timer
[Timer]
OnBootSec=420
OnCalendar=*:0,10,20,30,40,50
Unit=osrt-pkglistgen@.service

[Install]
WantedBy=timers.target
